# etc/systemd/system/steamos-desktop@.service
[Unit]
# =============================================================================
# SteamMachine-DIY - Desktop Session Service (KDE Plasma)
# =============================================================================
Description=SteamOS-DIY Desktop Session for %i
After=network.target systemd-user-sessions.service getty@tty1.service
# Conflicts ensures that starting Desktop Mode automatically stops Game Mode
Conflicts=steamos-gamemode@%i.service getty@tty1.service

[Service]
PAMName=login
User=%i
Group=%i
WorkingDirectory=/home/%i
UtmpIdentifier=tty1
UtmpMode=user

# Load Global Master Configuration (Agnostic Paths/Tags)
EnvironmentFile=/etc/default/steamos-diy

# TTY1 & Log Management
StandardInput=tty
# Redirecting desktop logs to the central diagnostic log
StandardOutput=append:/var/log/steamos-diy.log
StandardError=append:/var/log/steamos-diy.log
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes

# DYNAMIC ENVIRONMENT - Overriding Game Mode settings for a full KDE Session
Environment=XDG_SESSION_TYPE=wayland
Environment=XDG_CURRENT_DESKTOP=KDE
Environment=KDE_FULL_SESSION=true
Environment=QT_QPA_PLATFORM=wayland
# Disable specialized portals to use standard Desktop behavior
Environment=GTK_USE_PORTAL=0

# PRE-LAUNCH CLEANUP
# 1. Clear TTY1 for visual transition
ExecStartPre=-/usr/bin/bash -c "clear > /dev/tty1"
# 2. Force-release GPU (DRI) from any lingering Gamescope/Steam processes
ExecStartPre=-/bin/sh -c 'if /usr/bin/fuser /dev/dri/card* >/dev/null 2>&1; then /usr/bin/fuser -k /dev/dri/card*; fi'
# 3. Aggressive cleanup of stale Wayland sockets to prevent startup failures
ExecStartPre=-/usr/bin/bash -c "/usr/bin/rm -f /run/user/$(id -u %i)/wayland-0*"
# 4. Trigger the Desktop Mode Splash Screen
ExecStartPre=-/usr/local/bin/steamos-helpers/steamos-splash desktop

# LAUNCH SEQUENCE
# Start the native KDE Plasma Wayland session
ExecStart=/usr/bin/startplasma-wayland

# SESSION MANAGEMENT
Restart=no
KillMode=mixed
KillSignal=SIGTERM
# Fast timeout for seamless switching
TimeoutStopSec=500ms

[Install]
# Hook into the graphical target sequence
WantedBy=graphical.target
